# src/secrets/examples/systemd.service
# Example systemd service file for running with vault secrets
#
# The VAULT_PASSWORD should be set via systemd's LoadCredential or
# EnvironmentFile pointing to a file only readable by root.
#
# Copy to: /etc/systemd/system/aot-server.service
# Then: sudo systemctl daemon-reload && sudo systemctl start aot-server

[Unit]
Description=Agents of Treachery Game Server
After=network.target

[Service]
Type=simple
User=aot
Group=aot
WorkingDirectory=/opt/aot

# Option 1: Load password from a secure file (recommended)
# Create /etc/aot/vault-password with 0400 permissions owned by root
# The file should contain just the password, no newline
EnvironmentFile=/etc/aot/vault-password

# Option 2: Use systemd credentials (requires systemd 250+)
# LoadCredential=vault-password:/etc/aot/vault-password
# Environment=VAULT_PASSWORD=%d/vault-password

# Application environment
Environment=NODE_ENV=production
Environment=VAULT_ENV=prod
Environment=LLM_PROVIDER=anthropic
Environment=PORT=3001

# Start command - secrets are loaded at runtime
ExecStart=/usr/bin/node dist/server/index.js

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

# Restart policy
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
